# mlocate.timer
[Unit]
Description=locate/updatedb timer

[Timer]
OnCalendar=*-*-* 23:00:00

